# systemd/.config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStartPost=/usr/bin/emacsclient -e '(mu4e-update-index)'
