# easyeffecs/.config/systemd/user/easyeffects.service
[Unit]
Description=Easyeffects

[Service]
Type=simple
ExecStart=easyeffects --gapplication-service
Restart=always
RestartSec=1

[Install]
WantedBy=default.target